# miner.timer
[Unit]
Description=start miner.service after delay 3 min

[Timer]
OnBootSec=3min


[Install]
WantedBy=multi-user.target

