# miner.timer
[Unit]
Description=start miner.service after delay 5 min

[Timer]
OnBootSec=5min


[Install]
WantedBy=multi-user.target

